# linux/debian-package/DEBIAN/prerm
#!/bin/bash
################################################################################
# prerm - Pre-removal script for openpath-dnsmasq package
# Stops all services before package removal
################################################################################

set -e

case "$1" in
    remove|upgrade|deconfigure)
        echo "Deteniendo servicios openpath..."
        
        # Stop timers first
        systemctl stop openpath-dnsmasq.timer 2>/dev/null || true
        systemctl stop dnsmasq-watchdog.timer 2>/dev/null || true
        
        # Stop services
        systemctl stop captive-portal-detector.service 2>/dev/null || true
        systemctl stop dnsmasq 2>/dev/null || true
        
        # Disable services (but don't fail if they don't exist)
        systemctl disable openpath-dnsmasq.timer 2>/dev/null || true
        systemctl disable dnsmasq-watchdog.timer 2>/dev/null || true
        systemctl disable captive-portal-detector.service 2>/dev/null || true
        
        # Wait for dnsmasq to stop and kill residual processes
        sleep 2
        if pgrep -x dnsmasq >/dev/null 2>&1; then
            pkill -9 dnsmasq 2>/dev/null || true
            sleep 1
        fi
        
        echo "✓ Servicios detenidos"
        ;;
    
    failed-upgrade)
        ;;
    
    *)
        echo "prerm called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
