# linux/debian-package/lib/systemd/system/openpath-dnsmasq.service
[Unit]
Description=Update OpenPath DNS
After=network-online.target dnsmasq.service
Wants=network-online.target
Requires=dnsmasq.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openpath-update.sh
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
